# infra/supervisor/myblumbot.conf
[program:myblumbot]
autorestart=true
autostart=true
command=/root/myblumbot/scripts/prod_start.sh
stdout_logfile=/root/myblumbot/logs/out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/root/myblumbot/logs/err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
stopasgroup=true